# root/etc/s6-overlay/s6-rc.d/init-mod-deluge-yarss-plugin-add-package/up
/etc/s6-overlay/s6-rc.d/init-mod-deluge-yarss-plugin-add-package/run